# com.fifafantasy.datasync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd user agent for periodic FIFA Fantasy data sync.

  This agent periodically fetches match data from the remote API and updates
  sampresp.json. The cadence is time-of-day aware: every 30 minutes during the
  22:00-11:00 IST match window, every 3 hours otherwise (decided inside
  sync_data.py). If the API is down, it silently fails and will retry on the
  next interval. The publish workflow reads from sampresp.json, so an API
  outage never blocks publishing.

  Install:
    cp com.fifafantasy.datasync.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.fifafantasy.datasync.plist

  Check status:
    launchctl list | grep datasync
    tail -f ~/Library/Logs/fifafantasy-datasync.log

  Force a manual run (for testing):
    launchctl kickstart -k gui/$UID/com.fifafantasy.datasync

  NOTE: The project lives in ~/TProjects/ (not TCC-protected), so launchd can
  access it directly without needing Full Disk Access grants.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fifafantasy.datasync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/shishirmathur/TProjects/fifa_fantasy/sync_data.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shishirmathur/TProjects/fifa_fantasy</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Tick every 1 hour (3600 seconds). The script (sync_data.py) decides whether
         to actually fetch based on a time-of-day window (every 30 min during the
         22:00-11:00 IST match window, every 3 hours otherwise) and re-arms the
         matching pmset wake. The effective cadence is limited by launchd ticks
         (1h), so hot-window syncs will happen at most once per hour while awake. -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Log output (must live outside ~/Downloads, ~/Documents, ~/Desktop due
         to TCC protection; launchd cannot create files in those locations). -->
    <key>StandardOutPath</key>
    <string>/Users/shishirmathur/Library/Logs/fifafantasy-datasync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/shishirmathur/Library/Logs/fifafantasy-datasync.log</string>
</dict>
</plist>
